# centos/roles/nagios/files/fastcgi-wrapper.conf
start on runlevel [2345]
stop on runlevel [016]
respawn

exec /usr/bin/perl /usr/local/bin/fastcgi-wrapper.pl
