# debian/roles/hl/files/hl.conf
start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]
respawn

console none

setuid nobody
setgid nogroup

env LANG=en_US.UTF-8

script
  cd /srv/web/hl/
  ./.cabal-sandbox/bin/hl
end script